# share/docker/etc/services.d/02-thin/finish
#!/bin/sh -e

echo "[services.d] stopping service: thin"
