# share/docker/etc/services.d/02-thin/run
#!/bin/sh -e

echo "[services.d] starting service: thin"

su - redmine -c "/usr/bin/thin start -C config/thin.yml" > /dev/null 2>&1
